# arp_capture.service
[Unit]
Description=ARP packet capture service
After=network.target

[Service]
ExecStart=/etc/arp_capture/arp_capture.sh
Restart=always
User=root
Group=root
Environment=PATH=/usr/bin:/usr/local/bin

[Install]
WantedBy=multi-user.target
